# gluu-agent.service
[Unit]
Description=Gluu Agent Daemon

[Service]
Type=simple
ExecStart=/usr/bin/gluu-agent

[Install]
WantedBy=multi-user.target
