# gluu-agent.service
[Unit]
Description=Gluu Agent One Time Recover

[Service]
Type=oneshot
ExecStart="/usr/bin/gluu-agent recover"

[Install]
WantedBy=multi-user.target
